# conf/supervisor-app.conf
[program:app-uwsgi]
command = uwsgi --ini /home/docker/code/uwsgi.ini

[program:nginx-app]
command=/usr/sbin/nginx
stdout_events_enabled=true
stderr_events_enabled=true

