# conf/supervisor-app.conf
[program:app-uwsgi]
command = uwsgi --ini /home/docker/code/conf/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
